# deploy/aragora-backup.service
[Unit]
Description=Aragora Database Backup
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aragora
Environment="PATH=/home/ubuntu/aragora/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/ubuntu/aragora/venv/bin/python /home/ubuntu/aragora/scripts/backup_databases.py --keep 7

[Install]
WantedBy=multi-user.target
